# flightwall.service
[Unit]
Description=Flight Tracker LED Matrix Dashboard
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/justin/flightwall
ExecStart=/usr/bin/python3 app.py
Restart=always
RestartSec=10

# Environment (optional: set API keys, HOME_LAT, HOME_LON, etc. here or in .env)
# EnvironmentFile=/home/justin/flightwall/.env

[Install]
WantedBy=multi-user.target
