# flightwall.service
[Unit]
Description=Ribs FlightWall LED Tracker Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
# Adjust this path based on where the repository is cloned on the Pi
WorkingDirectory=/home/pi/flightwall
# Ensure nmcli scripts run before the Python app
ExecStartPre=/bin/bash /home/pi/flightwall/wifi_manager.sh
# Requires virtual environment or globally installed packages (e.g. via apt)
ExecStart=/usr/bin/python3 /home/pi/flightwall/main.py

Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=flightwall

[Install]
WantedBy=multi-user.target